# docs/kmscon.service.in
[Unit]
Description=KMS System Console
Documentation=man:kmscon(1)
After=systemd-user-sessions.service
After=rc-local.service

[Service]
ExecStart=@bindir@/kmscon -l /bin/login

[Install]
WantedBy=multi-user.target
